# schedule.crontab
53 18 * * * * python3 /Users/azuev/projects/weekly_report/schedule.py
